# adblocker.service
[Unit]
Description=Local DNS-based ad blocker
After=network.target

[Service]
ExecStart=/usr/local/bin/adblocker run
Restart=on-failure
RestartSec=5

# Data directory for blocklists
Environment=ADBLOCKER_DATA=/var/lib/adblocker

# Allow binding port 53 without full root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

# Minimal sandbox
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/adblocker
PrivateTmp=true

[Install]
WantedBy=multi-user.target
